# roles/maintenance_lock/templates/maintenance-lockdown.service.j2
[Unit]
Description=Maintenance lockdown (nologin + kick users)
After=network.target
ConditionPathExists=!/etc/nologin

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/maintenance-lockdown.sh
